# board/beaglebone/cron.d/cron-check-wifi
*/3 * * * * /etc/cron.d/check_wifi_status.sh